# scripts/codedeploy/templates/qmr-otel-refresh.timer
[Unit]
Description=Poll SSM every 5 min for QMR OTel Collector config/env changes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
