# scripts/codedeploy/templates/qmr-worker.service
[Unit]
Description=QMR Worker Service
After=network.target

[Service]
Type=simple
User=${SERVICE_USER}
Group=${SERVICE_GROUP}
WorkingDirectory=${APP_DIR}

# Python and application paths
Environment=PATH=${APP_DIR}/.venv/bin:/usr/local/bin:/usr/bin:/bin

# Model caching on instance storage (NVMe SSD)
Environment=TORCH_HOME=${MODEL_CACHE_DIR}
Environment=HF_HOME=${MODEL_CACHE_DIR}/huggingface
Environment=TORCHINDUCTOR_CACHE_DIR=${MODEL_CACHE_DIR}/inductor
Environment=TRITON_CACHE_DIR=${MODEL_CACHE_DIR}/triton
Environment=XDG_CACHE_HOME=${MODEL_CACHE_DIR}/xdg-cache

# Worker storage (image/template/result/cache dirs). Must be set at process
# spawn — worker/paths.py reads it at module-import time, before SSM loads.
Environment=WORKER_STORAGE_DIR=${WORKER_STORAGE_DIR}

# SSM Parameter Store path (Python worker loads these at startup)
Environment=SSM_PARAMETER_PATH=${SSM_PARAMETER_PATH}

ExecStart=${APP_DIR}/.venv/bin/python -m worker.main

# Restart policy
Restart=always
RestartSec=5

# Logging to file (for CloudWatch Agent) and journald
StandardOutput=append:/var/log/qmr-worker.log
StandardError=append:/var/log/qmr-worker.log
SyslogIdentifier=${SERVICE_NAME}

[Install]
WantedBy=multi-user.target
